# remap-logitech.service
[Service]
WorkingDirectory=/
ExecStart=bash -c "evremap remap /home/o/arch/evremap/logitech.toml"
Restart=always

[Install]
WantedBy=NetworkManager.service
